# webapp/devops/roles/supervisor/templates/supervisor.apps.ini
[program:{{ item.name }}]
command={{ item.command }}
numprocs=1
directory={{ item.directory }}
stdout_logfile=/var/log/{{ item.name }}.log
autostart=true
autorestart=true
user={{ user }}
